# focusd.service
[Unit]
Description=Updates /etc/hosts based on active user group
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/focusd.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target